# deploy/gunicorn.conf
description "Gunicorn application server handling solo"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data
chdir /home/ubuntu/solo

exec gunicorn --workers 4 --bind unix:/home/ubuntu/solo/solo.sock solo.wsgi:application